# rootfs/etc/services.d/freshclam/run
#!/usr/bin/execlineb -P
/usr/bin/freshclam --foreground --checks 6