# rootfs/etc/services.d/clamd/run
#!/usr/bin/execlineb -P
/usr/sbin/clamd --foreground